# bin/bridge/assets/reth-mainnet.service
[Unit]
Description=Scroll Mainnet - Reth
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=ubuntu
Group=ubuntu
UMask=0002
EnvironmentFile=/etc/bridge.env
ExecStart=ETHERSCAN_API_KEY=$ETHERSCAN_API_KEY /usr/local/bin/bridge --chain scroll-mainnet --metrics 0.0.0.0:9001 --debug.etherscan https://api.scrollscan.com/api

[Install]
WantedBy=multi-user.target